# crontab
# Run hourly quartz rebuilds
0 * * * * cd /app && npx quartz build

